# docker/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:cron]
command=/usr/sbin/cron -f

[program:node]
command=/usr/local/bin/forever /src/app/server/server.js
redirect_stderr=true
